# docker/php/crontab

# запуск команды обновления курсов раз в час через команду

0 * * * *      php /var/www/html/bin/console app:get-rate